# templates/phc2sys.service.j2
[Unit]
Description=Sync PHC to System clock
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/phc2sys -s {{ phc2sys_source }} -c {{ phc2sys_clock }} --step_threshold=1 --transportSpecific=1 {{ phc2sys_extra_opts }}
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target